# debian/DEBIAN/postinst
#!/bin/bash
set -e

# 在安装/升级后运行

APP_NAME="ems-simulate"
INSTALL_DIR="/usr/share/${APP_NAME}"
BIN_DIR="/usr/bin"
SERVICE_NAME="${APP_NAME}.service"

case "$1" in
    configure)
        # 1. 重新加载 systemd
        if command -v systemctl >/dev/null; then
            systemctl daemon-reload
            # 可选：安装后自动启动
            systemctl enable "${SERVICE_NAME}"
            systemctl start "${SERVICE_NAME}"
        fi
        
        # 打印服务管理提示
        GREEN='\033[0;32m'
        BLUE='\033[0;34m'
        YELLOW='\033[1;33m'
        NC='\033[0m'
        
        echo -e "${GREEN}✅ ${APP_NAME} 安装成功！${NC}"
        echo -e "${BLUE}服务管理命令:${NC}"
        echo -e "  启动服务: ${GREEN}sudo systemctl start ${SERVICE_NAME}${NC}"
        echo -e "  停止服务: ${YELLOW}sudo systemctl stop ${SERVICE_NAME}${NC}"
        echo -e "  查看状态: ${BLUE}sudo systemctl status ${SERVICE_NAME}${NC}"
        echo -e "  开机自启: ${BLUE}sudo systemctl enable ${SERVICE_NAME}${NC}"
        echo -e "  Web界面: http://127.0.0.1:8991"
        
        echo -e "\n${BLUE}关于:${NC}"
        echo -e "  作者: CDY"
        echo -e "  Email: 1755696012@qq.com"
        echo -e "  GitHub: https://github.com/600888/ems_simulate"
        echo -e "  文档: https://600888.github.io/ems_simulate/"

        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

exit 0
